# systemd/trading-backup.timer
[Unit]
Description=Trading Intelligence - Daily Backup Timer

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
